# c/crontab2/pgdpcanada.cron
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.ITtwhh/crontab installed on Sun Aug  9 06:26:53 2009)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# multiple times per hour:

10,25,40,55 1-23 * * * /c/tools/project_manager/automodify.php 
   25,40       0 * * * /c/tools/project_manager/automodify.php 

----------------------

03 * * * * /sharehome/htdocs/c/crontab2/update_user_counts.php 
01  0 * * * /sharehome/htdocs/c/crontab2/take_tally_snapshots.php
02  0 * * * /sharehome/htdocs/c/crontab2/log_project_states.php
05  0 * * * /sharehome/htdocs/c/crontab2/update_user_birthdays.php
10  0 * * * /sharehome/htdocs/c/crontab2/archive_projects.php
50  0 * * * find /home/pgdpcanada/htdocs/d/sp_check/* -mtime +2 -exec rm -f {} \;
51  0 * * * find /home/pgdpcanada/htdocs/d/automodify_logs/* -mtime +2 -exec rm -f {} \;
54  0 * * * /sharehome/htdocs/c/crontab2/extend_site_tally_goals.php
00  4 * * * /sharehome/htdocs/c/crontab2/onoff_special_event_queues.php
57  0 * * * mysqldump --protocol=TCP --host=localhost --user=usrDPCanada --password=CanadaWin10 --opt --all-databases > /home/pgdpcanada/sqlbackup/alldbbkp.sql

06  2 * * * sudo freshclam
08  2 * * * sudo clamscan -ir --max-filesize=10000000 / | mail -s "Virus scan report" tfsadmin
04 4 7,14,21,28 * * /sharehome/htdocs/c/crontab2/update_ppers.php
20 0 1 * * /sharehome/htdocs/c/crontab2/inactivate_users.php
35 0 1 * * /sharehome/htdocs/c/crontab2/notify_users.php
49 0 1 * * /sharehome/htdocs/c/crontab2/notify_old_pp.php


